# deploy/docker/supervisord.conf
[supervisord]
nodaemon=true
# Use the overwatch data directory
logfile=/opt/overwatch/data/supervisord.log
childlogdir=/opt/overwatch/data/
#logfile=/Users/re239/code/alice/overwatch/data/supervisord.log
#childlogdir=/Users/re239/code/alice/overwatch/data/
# 5 MB log file with 10 backup files
logfile_maxbytes=5000000
logfile_backups=10

[unix_http_server]
# (the path to the socket file)
file = /tmp/supervisor.sock
# socket file mode (default 0700)
chmod = 0700

# the below section must remain in the config file for RPC
# (supervisorctl/web interface) to work, additional interfaces may be
# added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
# use a unix:// URL  for a unix socket
serverurl=unix:///tmp/supervisor.sock

#[program:nginx]
#command=/usr/sbin/nginx
#redirect_stderr=true
## 5 MB log file with 10 backup files
#stdout_logfile_maxbytes=500000
#stdout_logfile_backups=10
#
#[program:overwatch]
#command=/bin/bash -c "/opt/overwatch/deploy/initOVERWATCH.sh docker"
#redirect_stderr=true
## 5 MB log file with 10 backup files
#stdout_logfile_maxbytes=500000
#stdout_logfile_backups=10

[program:startOverwatch]
command=overwatchDeploy -e config
#command=python /opt/overwatch/deploy/start.py -e config
#command=python /Users/re239/code/alice/overwatch/deploy/start.py -c /Users/re239/code/alice/overwatch/deploy/deployConfig.yaml
autorestart=false
startsecs=0
redirect_stderr=true
# 5 MB log file with 10 backup files
stdout_logfile_maxbytes=500000
stdout_logfile_backups=10

